# package-server/supervisord.conf
[program:package-listing]
user=apply2
command=/usr/bin/java -jar /home/arjun/package-listing.jar
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
